# crontab
Add in your crontab
0 0 * * * /path/to/your/script/script.sh
